# deploy/pyems-time-sync.service
# PyEMS controller-clock helper — called every minute by pyems-time-schedule.
# It contacts NTP only when the schedule in the PyEMS logs directory matches.
[Unit]
Description=PyEMS scheduled controller time synchronization
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/home/pi/PyEMS/.venv/bin/pyems-time-helper --state /home/pi/PyEMS/logs/time-settings.json --status /home/pi/PyEMS/logs/time-sync-status.json --sync-if-due
