# deploy/pyems-time-apply.service
# PyEMS controller-clock helper — apply a manual time or NTP schedule setting.
# install.sh generates this with the real checkout and service-user paths.
[Unit]
Description=PyEMS apply controller time settings

[Service]
Type=oneshot
User=root
ExecStart=/home/pi/PyEMS/.venv/bin/pyems-time-helper --apply
